# scripts/com.0xsloane.crypto-factor-screener-refresh.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.0xsloane.crypto-factor-screener-refresh</string>

  <key>WorkingDirectory</key>
  <string>/Users/asad.hussain/Crypto Factor Model</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/asad.hussain/Crypto Factor Model/.venv/bin/python</string>
    <string>/Users/asad.hussain/Crypto Factor Model/scripts/refresh_screener_cache.py</string>
    <string>--start</string>
    <string>2024-06-01</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>8</integer>
    <key>Minute</key>
    <integer>10</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/asad.hussain/Crypto Factor Model/logs/screener_refresh.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/asad.hussain/Crypto Factor Model/logs/screener_refresh.err.log</string>
</dict>
</plist>
